# nweb.service
[Unit]
Description=A simple static file webserver

[Service]
ExecStart=/usr/bin/nweb 80 /srv/nweb
Type=forking

[Install]
WantedBy=network.target
